# dixet/ansible/files/systemd/default.target.wants/kb-it.timer
[Unit]
Description=Knowledge Base: IT
PartOf=kb.target

[Timer]
OnBootSec=5min
OnCalendar=*:0/5
Unit=kb-it.service

[Install]
WantedBy=default.target