# dixet/ansible/files/systemd/kb-languages.timer
[Unit]
Description=Knowledge Base: Languages

[Timer]
OnBootSec=5min
OnCalendar=*:0/5
Unit=kb-languages.service

[Install]
WantedBy=default.target
